# contrib/systemd/hister.service
[Unit]
Description=Hister web history service
After=network.target

[Service]
Type=simple
User=hister
Group=hister
StateDirectory=hister
Environment=HISTER_DATA_DIR=/var/lib/hister
EnvironmentFile=-/etc/hister/hister.env
ExecStart=/usr/local/bin/hister listen
Restart=on-failure

# For privileged ports (< 1024), set both values to CAP_NET_BIND_SERVICE.
AmbientCapabilities=
CapabilityBoundingSet=

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
LockPersonality=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native
SystemCallFilter=@system-service ~@privileged
MemoryDenyWriteExecute=true
UMask=0077

[Install]
WantedBy=multi-user.target
